# systemd/containers-image-prune.timer
[Unit]
Description=Prune unused images weekly

[Timer]
OnCalendar=Mon *-*-* 23:00

[Install]
WantedBy=timers.target

