# systemd/fetch_news.timer
[Unit]
Description=Daily news fetching

[Timer]
OnCalendar=*-*-* 18:00:00

[Install]
WantedBy=timers.target