# ansible/roles/k3s/files/k3s-upgrade.timer
[Unit]
Description=K3s Upgrade Timer

[Timer]
OnCalendar=Mon *-*-* 3:00:00

[Install]
WantedBy=timers.target