# ansible/roles/k3s/files/k3s-upgrade.timer
[Unit]
Description=K3s Upgrade Timer

[Timer]
OnCalendar=*-*-* 4 * * 1

[Install]
WantedBy=timers.target