# ansible/roles/k3s/files/k3s-upgrade.service
[Unit]
Description=K3s Upgrade Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/k3s-upgrade.sh
